# deploy/systemd/pib4ecu-net-mode-apply.service
[Unit]
Description=Apply piB4ECU network mode policy at boot
# Run after the generic network stack is up so client Wi-Fi (NetworkManager) is not
# left in a broken state; car/hotspot still switches stack on the next lines.
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'cd /home/joe/Code/piB4ECU && ./scripts/pib4ecu-net-mode.sh apply'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
